# linux-customizations/random-wallpaper/random-wallpaper.timer
[Unit]
Description=Run random wallpaper hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
